# docker/supervisord.conf
; Supervisor program definitions for the combined API + UI container.
; Installed to /etc/supervisor/conf.d/app.conf; the base supervisord.conf
; includes everything in conf.d/*.conf. Started with `supervisord -n`.

[program:api]
command=uvicorn hf_readmit.api.app:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:streamlit]
command=streamlit run src/hf_readmit/ui/app.py --server.port 8501 --server.address 0.0.0.0
directory=/app
environment=HF_API_URL="http://localhost:8000"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
